# rasponics.service
[Unit]
Description=Aquaponics Automation Controller

[Service]
ExecStart=/home/hburke/rasponics/rasponics.py
Restart=always
